# container_hosts.cron
*/2 * * * * root /usr/local/bin/container_hosts
